# scripts/schoolflow-backup.timer
[Unit]
Description=SchoolFlow Pro Daily Backup Timer
Requires=schoolflow-backup.service

[Timer]
OnCalendar=daily
OnCalendar=02:00
# Ensure backup runs even if the machine was off
Persistent=true

[Install]
WantedBy=timers.target
