# scripts/schoolflow-backup.timer
[Unit]
Description=SchoolFlow Pro Daily Backup Timer

[Timer]
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=10m
Persistent=true
Unit=schoolflow-backup.service

[Install]
WantedBy=timers.target
